# supabase/systemd/hyper-auth-supabase-restore-drill.timer
[Unit]
Description=Daily trigger for due Hyper Authenticator Supabase restore drill

[Timer]
OnCalendar=*-*-* 04:30:00 UTC
Persistent=true
RandomizedDelaySec=15m
Unit=hyper-auth-supabase-restore-drill.service

[Install]
WantedBy=timers.target
